# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:snapshotter]
command=/opt/cassandra/bin/snapshot-scheduler.sh


[program:cassandra]
command=/opt/cassandra/bin/stups-cassandra.sh
